# vix.conf
[program:vix]
command=python3 main.py
directory=/app/lib
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/vix.log
stderr_logfile=/var/log/supervisor/vix.err
