# vix_service.conf
[program:vix_service]
command=/app/service/venv/bin/python3 -m visionatrix run --ui --flows_dir=/nc_app_vix_data/vix_flows --models_dir=/nc_app_vix_data/vix_models --tasks_files_dir=/nc_app_vix_data/vix_tasks_files --disable-smart-memory
directory=/app/service
autostart=true
autorestart=true
startsecs=30
stdout_logfile=/var/log/supervisor/vix_service.log
stderr_logfile=/var/log/supervisor/vix_service.err
